# dexboost.service
[Unit]
Description=DexBoost Backend Service
After=network.target

[Service]
Type=simple
User=nobody
Group=nobody
WorkingDirectory=/home/dexboost
ExecStart=/usr/bin/npm start
Restart=always
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=CONFIG_FILE=/home/dexboost/src/config.prod.ts

[Install]
WantedBy=multi-user.target 